# dhl-tracking-sync/deploy/dhl-tracking-sync.timer
[Unit]
Description=Run DHL tracking sync once daily (Rate-Limit 250/Tag beachten)

[Timer]
OnCalendar=*-*-* 06:30:00
Persistent=true

[Install]
WantedBy=timers.target
